# cm15d/conf/supervisor.conf.example
[program:cm15d]
command=/usr/local/bin/cm15d.py
user=cm15d
